# host/src/pack/dpkg/DEBIAN/postinst
#!/bin/bash

udevadm control --reload-rules
udevadm trigger

systemctl enable mwc
systemctl start mwc.service
