# host/src/pack/dpkg/DEBIAN/prerm
#!/bin/bash

systemctl stop mwc.service
systemctl disable mwc
